# deploy/streamnews-analyzer.service
[Unit]
Description=StreamNews Analyzer (FastAPI)
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=streamnews
Group=streamnews
WorkingDirectory=/opt/streamnews/analyzer
EnvironmentFile=/opt/streamnews/.env
Environment=PATH=/opt/streamnews/.venv/bin:/usr/bin
ExecStart=/opt/streamnews/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
